# etc/supervisor/conf.d/cron.conf
[program:cron]
command = /usr/sbin/cron -f -L 15
autorestart = true
stdout_logfile = /var/log/supervisor/cron.log 
stderr_logfile = /var/log/supervisor/cron.err 
